# launchd/com.greg.shout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.greg.shout</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{SHOUT_BIN}}</string>
    <string>daemon</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ProcessType</key>
  <string>Interactive</string>

  <key>StandardOutPath</key>
  <string>/tmp/shout.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/shout.err.log</string>

  <!-- Throttle restart attempts so a crash loop does not eat CPU. -->
  <key>ThrottleInterval</key>
  <integer>5</integer>
</dict>
</plist>
